# CheckIpTask.cron
*/5 * * * * sh /home/wangzhe//DynamicDNS/UpdateIP.sh
